# example/files/supervisor/apache.conf
; supervisor configuration for Apache
[program:apache]
command = /bin/sh -c '. /etc/apache2/envvars && exec /usr/sbin/apache2 -D FOREGROUND'
